# infra/seo-cron/kidsbayarea-seo-cron.timer
[Unit]
# Weekly (Thu). 1x/week because GSC needs ~10-14 days to reflect a snippet
# change — running faster than the measurement window confounds attribution
# (see playbook.md §"single-page cooldown"). Thu is the free day on this shared
# box: Mon=visacub(seo-weekly), Tue=citationmap(seo-cron), Wed=pickfromvideo
# (seo-cron-pfv), Sat=eurkart(seo-cron-eurkart). 14:13 clears their 13:13-13:22
# cluster too.
Description=Run kidsbayarea-seo-cron.service weekly (Thu 14:13 UTC; the free day/slot on the shared box)

[Timer]
OnCalendar=Thu *-*-* 14:13:00 UTC
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
