# infra/seo-cron/kidsbayarea-seo-cron-failure-notify.service
[Unit]
Description=Email alert when kidsbayarea-seo-cron.service fails (wired via OnFailure=)

[Service]
Type=oneshot
User=ec2-user
Environment=HOME=/home/ec2-user
ExecStart=/bin/bash -c 'aws ses send-email --region us-east-1 \
  --from seo-cron@anureka.com --to tom@anureka.com \
  --subject "[SEO cron] FAILED — kidsbayarea.com $(date -u +%%F)" \
  --text "kidsbayarea-seo-cron.service exited with failure on visacub-seo-cron (i-0fee473f9b347c706). Last journal lines:\n\n$(journalctl -u kidsbayarea-seo-cron.service -n 30 --no-pager -o cat | tail -c 3000)"'
