# infra/seo-cron/seo-cron-kidsbayarea.service
[Unit]
Description=Weekly SEO optimization run for kidsbayarea.com (GSC+Bing+DataForSEO -> Claude Code on Bedrock -> PR + report)
Wants=network-online.target
After=network-online.target
OnFailure=seo-cron-kidsbayarea-failure-notify.service

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/opt/seo-cron-kidsbayarea
ExecStart=/opt/seo-cron-kidsbayarea/kidsbayarea/infra/seo-cron/run.sh
TimeoutStartSec=2h
Environment=HOME=/home/ec2-user
